# Initializing_Wifi/ip_boot.service
[Unit] 
After=network-online.target
Description=Get IP Address

[Service]
ExecStart=/home/pi/hello.sh

[Install]
WantedBy=multi-user.target

